# service/gateway/run
#!/bin/sh
exec 2>&1
exec su - user -c "cd /home/user/mozilla-iot/gateway && ./run-app.sh"
