# service/gateway/run
#!/bin/sh

export MOZIOT_HOME=/home/node/.mozilla-iot
chown -R node:node "${MOZIOT_HOME}"
exec 2>&1
exec su - node -c "cd /home/node/mozilla-iot/gateway && ./run-app.sh"
