# service/intent-parser/run
#!/bin/sh

exec 2>&1
exec su - node -c "cd /home/node/mozilla-iot/intent-parser && ./intent-parser-server.py"
