# cron.d/renew-certs
0 */12 * * * root /opt/scripts/renew-certs
